# templates/sensu-web.service.erb
[Unit]
Description=Sensu Go Web
After=network-online.target multi-user.target
Wants=network-online.target

[Service]
Environment=NODE_ENV=production
Environment=PORT=<%= @port %>
Environment=API_URL=<%= @api_url %>
WorkingDirectory=<%= @install_dir %>
User=<%= @service_user %>
Group=<%= @service_group %>
ExecStart=/usr/bin/yarn node scripts serve

[Install]
WantedBy=multi-user.target
